# oracle-jdk-17-lts/etc/profile.d/jdk-17.sh
export JAVA_HOME=/usr/lib/jvm/jdk-17/
export PATH=$JAVA_HOME/bin:$PATH
